# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/sbin/php5-fpm -F
autostart=true
autorestart=unexpected

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
